# telegram-bridge.service
[Unit]
Description=Telegram Bridge Worker
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/telegram-bridge
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/telegram-bridge/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bridge

[Install]
WantedBy=multi-user.target
