# telegram-bridge-bot.service
[Unit]
Description=Telegram Bridge Management Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/telegram-bridge
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/telegram-bridge/venv/bin/python src/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bridge-bot

[Install]
WantedBy=multi-user.target
